# scripts/kalinka.service
[Unit]
Description=Run kalinka backend
After=multi-user.target

[Service]
Type=simple
User=kalusr
WorkingDirectory=/opt/kalinka
ExecStart=/usr/bin/kalinka_server.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target